# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:c3-web]
command=/usr/bin/start_c3-web.sh